# mkosi.extra/install-files/settings/home_audioserver_.config/systemd/user/audioserver_restart_broadcast.timer
[Unit]
Description=Restart broadcast daily

[Timer]
OnCalendar=3:00
Persistent=false

[Install]
WantedBy=timers.target
