# mkosi.extra/install-files/settings/home_audioserver_.config/systemd/user/audioserver_restart_broadcast.service
[Unit]
Description=Restart broadcast

[Service]
Type=simple
ExecStart=/settings/audioserver_restart_broadcast.sh
