# install-files/settings/home_audioserver_.config/systemd/user/audioserver-user-update.service
[Unit]
Description=audioserver user update script
ConditionPathExistsGlob=/home/audioserver/user.version

[Service]
Type=oneshot
ExecStart=/install-files/install/user-update.sh

[Install]
WantedBy=multi-user.target
