# install-files/settings/home_audioserver_.config/systemd/user/audioserver_reload_mopidy.timer
[Unit]
Description=Reload mopidy database daily

[Timer]
OnCalendar=daily
Persistent=false    
 
[Install]
WantedBy=timers.target
